# profile.d/init_env.local
# This file is the initial environment set up for all shells, to help set up
# a few things, like where the dotfiles are located.
