# profile.d/languagerc.local
#!/bin/bash

# Local language defaults, e.g.:

# rvm use 2.3
# nvm use v6.2.0
# export PYENV_VERSION=3.5.2
