# profile.d/netrc.local
# Set up github api tokens (e.g. vim-rhubarb)
machine api.github.com login <user> password <token>
