# profile.d/plugins.sh
####################
# External plugins #
####################

# macOS plugins
if is_osx; then
    source_file $DOTFILES/profile.d/external/tab.sh # OSX new tab
fi
