# profile.d/plugins.sh
####################
# External plugins #
####################

source_file /usr/local/opt/autoenv/activate.sh # Autoenv
source_file $DOTFILES/profile.d/external/tab.sh # OSX new tab
