# profile.d/plugins.sh
####################
# External plugins #
####################

# Autoenv
source_file /usr/local/opt/autoenv/activate.sh
