# profile.d/profilerc.local
#!/bin/bash

# Local profile specific environment settings

# Homebrew
export HOMEBREW_GITHUB_API_TOKEN=

# Load local language configuration, if confirmed
[[ -r $HOME/.languagerc ]] &&  ask "Load languages?" N && . $HOME/.languagerc
