# profile.d/profilerc.local
#!/bin/bash

# Local profile specific environment settings

# Homebrew
export HOMEBREW_GITHUB_API_TOKEN=
