# profile.d/profilerc.local
#!/bin/bash

# Local profile configuration
# These settings are applied last and can be used to override any prior configuration

# Homebrew
export HOMEBREW_GITHUB_API_TOKEN=
